# airootfs/etc/profile.d/boxpath.sh
export PATH=$PATH:/opt/toybox:/opt/busybox:/opt/sbase/bin
export MANPATH=$MANPATH:/opt/sbase/man
